# test/agent/microshift-test-agent.service
[Unit]
Description=MicroShift Test Agent
Before=microshift.service sshd.service

[Service]
ExecStart=/usr/bin/microshift-test-agent.sh
# Default value of KillMode is control-group which results in
# sending SIGTERM to the process and its child processes.
# By changing to `mixed`, only the main process will get the
# SIGTERM and cleanup activities won't be interrupted.
KillMode=mixed

[Install]
WantedBy=multi-user.target
