# test/agent/microshift-test-agent.service
[Unit]
Description=MicroShift Test Agent
Before=microshift.service

[Service]
ExecStart=/bin/bash /usr/bin/microshift-test-agent.sh

[Install]
WantedBy=multi-user.target
