# packaging/systemd/microshift.service
[Unit]
Description=MicroShift
Wants=network-online.target crio.service
After=network-online.target crio.service

[Service]
WorkingDirectory=/usr/bin/
ExecStartPre=/bin/bash -c 'mv /var/log/microshift.log /var/log/$$(date +%%T)-microshift.log'
ExecStart=microshift run
StandardOutput=file:/var/log/microshift.log
Restart=always
User=root

[Install]
WantedBy=multi-user.target
