# ansible/roles/setup-microshift-host/templates/cadvisor.service.j2
[Unit]
Description=cAdvisor

[Service]
ExecStart=/usr/bin/cadvisor

[Install]
WantedBy=multi-user.target
