# ansible/roles/install-logging-exporters/templates/cadvisor.service.j2
[Unit]
Description=cAdvisor

[Service]
ExecStart=/usr/bin/cadvisor -port=8082

[Install]
WantedBy=multi-user.target
