# packaging/observability/microshift-observability.service
[Unit]
Description=MicroShift Observability
After=microshift.service

[Service]
ExecStart=/usr/bin/opentelemetry-collector --config=/etc/microshift/opentelemetry-collector.yaml
Restart=always
User=root
TimeoutStartSec=60

# It takes a bit for the certs to be created. This service will reach it's burst limit almost immediately, pretty much
# guaranteeing that it will reach the restart limit before it can possibly succeed.
RestartSec=200ms

[Install]
RequiredBy=multi-user.target