# packaging/observability/microshift-observability.service
[Unit]
Description=MicroShift Observability
After=microshift.service
Requires=microshift.service
ConditionPathExists=/var/lib/microshift/resources/observability-client/kubeconfig

[Service]
Environment=KUBECONFIG=/var/lib/microshift/resources/observability-client/kubeconfig
Environment=K8S_NODE_NAME="%l"
ExecStartPre=/usr/bin/mkdir -p /var/lib/microshift-observability
ExecStart=/usr/bin/opentelemetry-collector --config=/etc/microshift/observability/opentelemetry-collector.yaml
Restart=always
User=root

[Install]
WantedBy=microshift.service
